# configfiles/rc.local
#!/bin/sh -e

## AUTHOR:       AISK11
## LOCATION:     /etc/rc.local (0644 root:root)
## DEPENDENCIES: -
## DESCRIPTION:  Startup commands.

## Unblock Wi-Fi.
rfkill block all
rfkill unblock wlan

exit 0
